# harbour-batteryoverlay2.service.in
[Unit]
Description=BatteryOverlay
After=dbus.socket booster-generic.service lipstick.service
Requires=dbus.socket booster-generic.service lipstick.service

[Service]
Type=dbus
BusName=harbour.batteryoverlay.overlay
ExecStart=/usr/bin/invoker --type=generic -- @@SAILJAIL@@ /usr/bin/harbour-batteryoverlay2 -daemon

[Install]
WantedBy=user-session.target
